# build/macos/com.cloudatlas.unified-hifi-control.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cloudatlas.unified-hifi-control</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/unified-hifi-control</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/unified-hifi-control.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/unified-hifi-control.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>8088</string>
        <key>CONFIG_DIR</key>
        <string>/usr/local/var/unified-hifi-control</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/usr/local/var/unified-hifi-control</string>
</dict>
</plist>
